# lib/systemd/system/failsafe-graphical.target
[Unit]
Description=Graphical failsafe fallback
Requires=multi-user.target failsafe-x.service
After=multi-user.target
AllowIsolate=yes
