# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
What=/dev/disk/by-uuid/c486629a-8c47-4d0b-88aa-b6a9c08de0e7
Where=/
Type=ext4
Options=errors=remount-ro
