# var/lib/dpkg/info/nplan.list
/.
/etc
/etc/netplan
/lib
/lib/netplan
/lib/netplan/generate
/lib/systemd
/lib/systemd/system
/lib/systemd/system/netplan-wpa@.service
/lib/systemd/system-generators
/usr
/usr/sbin
/usr/sbin/netplan
/usr/share
/usr/share/doc
/usr/share/doc/netplan
/usr/share/doc/netplan/netplan.html
/usr/share/doc/nplan
/usr/share/doc/nplan/changelog.gz
/usr/share/doc/nplan/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/netplan.5.gz
/lib/systemd/system-generators/netplan
