# var/lib/dpkg/info/repowerd.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/com.canonical.Unity.Screen.conf
/etc/dbus-1/system.d/com.canonical.powerd.conf
/etc/dbus-1/system.d/com.canonical.repowerd.conf
/etc/init
/etc/init/repowerd.conf
/etc/logrotate.d
/etc/logrotate.d/repowerd
/etc/rsyslog.d
/etc/rsyslog.d/20-repowerd.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/repowerd.service
/usr
/usr/sbin
/usr/sbin/repowerd
/usr/sbin/repowerd-cli
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/repowerd.py
/usr/share/doc
/usr/share/doc/repowerd
/usr/share/doc/repowerd/changelog.Debian.gz
/usr/share/doc/repowerd/copyright
/usr/sbin/powerd-cli
