# run/user/1000/systemd/user/unity7.service.d/graphical-session-pre.conf
[Unit]
After=graphical-session-pre.target
