# lib/systemd/system/display-manager.service.d/xdiagnose.conf
# Fallback to xdiagnose if we couldn't start the display-manager

[Unit]
OnFailure=failsafe-graphical.target

